# slp-pkgmgr.service
[Unit]
Description=SLP Package Manager

[Service]
BusName=com.samsung.slp.pkgmgr
ExecStart=/usr/bin/pkgmgr-server