# slp-pkgmgr.service
[Unit]
Description=SLP Package Manager

[Service]
BusName=org.tizen.pkgmgr
ExecStart=/usr/bin/pkgmgr-server
